# kube/research-ai-ai.container
[Unit]
Description=Research AI LLM Service (GPU)
After=nvidia-persistenced.service
Wants=nvidia-persistenced.service

[Service]
Restart=always
RestartSec=5
StartLimitIntervalSec=0
StartLimitBurst=0

[Container]
ContainerName=research-ai-ai
Image=docker.io/ollama/ollama:latest
Network=research-ai-net

Environment=NVIDIA_VISIBLE_DEVICES=all
Environment=NVIDIA_DRIVER_CAPABILITIES=compute,utility
AddDevice=nvidia.com/gpu=all

# Expose to localhost on the server so the SSH tunnel can reach it
PublishPort=127.0.0.1:11434:11434 

# Keep models loaded indefinitely (no auto-unload after idle timeout).
Environment=OLLAMA_KEEP_ALIVE=-1

# Ollama model storage — uses a host directory so models persist across
# container recreates. ~/.ollama is the conventional default location.
Volume=/root/.ollama:/root/.ollama

[Install]
WantedBy=multi-user.target
