# kube/research-ai-ricgraph.container
[Unit]
Description=Research AI Ricgraph
After=research-ai-neo4j.service

[Service]
Restart=always
RestartSec=5
StartLimitIntervalSec=0
StartLimitBurst=0

[Container]
ContainerName=research-ai-ricgraph
Image=ghcr.io/utrechtuniversity/ricgraph:latest
Network=research-ai-net

# Expose to localhost on the server so the SSH tunnel can reach it
# 3030 is for Ric explorer
PublishPort=127.0.0.1:3030:3030

EnvironmentFile=/etc/research-ai/research-ai-prod.env

# Override ricgraph.ini to point at the external Neo4j container
Volume=/etc/research-ai/ricgraph.ini:/app/ricgraph/ricgraph.ini:ro

[Install]
WantedBy=multi-user.target
