# kube/research-ai-api.container
[Unit]
Description=Research AI API (GPU)
Requires=research-ai-frontend.service
After=research-ai-frontend.service

[Service]
Restart=always
RestartSec=5
StartLimitIntervalSec=0
StartLimitBurst=0

[Container]
ContainerName=research-ai-api
Image=research-ai-api:prod
Network=container:research-ai-frontend

EnvironmentFile=/etc/research-ai/research-ai.env
Environment=NVIDIA_VISIBLE_DEVICES=all
Environment=NVIDIA_DRIVER_CAPABILITIES=compute,utility
AddDevice=nvidia.com/gpu=all

[Install]
WantedBy=multi-user.target