# kube/research-ai-neo4j.container
[Unit]
Description=Research AI Neo4j Database
Requires=research-ai-frontend.service
After=research-ai-frontend.service

[Service]
Restart=always
RestartSec=5
StartLimitIntervalSec=0
StartLimitBurst=0

[Container]
ContainerName=research-ai-neo4j
Image=docker.io/library/neo4j:5
Network=container:research-ai-frontend

Volume=neo4j-data:/data

EnvironmentFile=/path/to/your/project.env

[Install]
WantedBy=multi-user.target